# 95-hdmi-plug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/baum/.Xauthority", RUN+="/home/baum/Documents/bash/togdisp.sh"

